# rootfs-overlay/climate_prodtest_base/usr/etc/udev/rules.d/80-iio.rules
ACTION=="add", SUBSYSTEM=="iio", ATTR{label}=="ClimateProdTest_DAC1", RUN+="/bin/sh -c 'chgrp -R torizon /sys%p; chmod -R 774 /sys%p'"
