# rootfs-overlay/climate_prodtest_base/usr/etc/udev/rules.d/99-label-printer.rules
ACTION=="add|change", SUBSYSTEM=="usbmisc", ATTRS{idVendor}=="0a5f", SYMLINK +="label_printer", TAG+="systemd", ENV{SYSTEMD_ALIAS}+="/sys/subsystem/usbmisc/label_printer", RUN+="/etc/udev/scripts/add-printer.sh"
